# examples/systemd/spbackup-list.service
[Unit]
Description=SharePoint List Backup (spbackup list)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=spbackup
Group=spbackup

# ── Paths — edit these ───────────────────────────────────────────────────
WorkingDirectory=/opt/spbackup
ExecStart=/usr/bin/pwsh /opt/spbackup/spbackup.ps1 list backup --site-url "https://contoso.sharepoint.com/sites/team" --list "Tasks" --out /var/lib/spbackup/lists --verbose

# ── Credentials — use a drop-in or systemd-creds ────────────────────────
# Option A: EnvironmentFile (chmod 600, owned by root or spbackup user)
EnvironmentFile=/etc/spbackup/env
# The env file should contain (NO quotes around values):
#   TENANT_ID=xxxxxxxx-xxxx-xxxx-xxxx-xxxxxxxxxxxx
#   CLIENT_ID=xxxxxxxx-xxxx-xxxx-xxxx-xxxxxxxxxxxx
#   CLIENT_SECRET=your-secret-value
#
# Option B: systemd credentials (systemd 250+)
# LoadCredentialEncrypted=client_secret:/etc/credstore/spbackup.secret
# Then read in a wrapper script.

# ── Hardening ────────────────────────────────────────────────────────────
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/spbackup
PrivateTmp=true

# ── Timeouts & Resources ────────────────────────────────────────────────
TimeoutStartSec=1800
MemoryMax=512M

[Install]
WantedBy=multi-user.target
