# examples/systemd/spbackup-loop.service
[Unit]
Description=SharePoint Loop Backup (spbackup loop)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=spbackup
Group=spbackup

# ── Paths — edit these ───────────────────────────────────────────────────
WorkingDirectory=/opt/spbackup
ExecStart=/usr/bin/pwsh /opt/spbackup/spbackup.ps1 loop backup --url "https://YOUR-ORG.sharepoint.com/:f:/s/SITE/XXXXX" --out /var/lib/spbackup/loop --verbose

# ── Credentials — use a drop-in or systemd-creds ────────────────────────
EnvironmentFile=/etc/spbackup/env

# ── Hardening ────────────────────────────────────────────────────────────
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/spbackup
PrivateTmp=true

# ── Timeouts & Resources ────────────────────────────────────────────────
TimeoutStartSec=1800
MemoryMax=512M

[Install]
WantedBy=multi-user.target
