# examples/windows/SPBackup-Library.xml
<?xml version="1.0" encoding="UTF-16"?>
<!--
  Windows Task Scheduler XML — SharePoint Document Library Backup
  Import:  schtasks /Create /TN "SPBackup-Library" /XML "SPBackup-Library.xml"
  Edit the paths, URL, library name, and credentials below before importing.
-->
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>Back up a SharePoint document library (all files with folder structure) every 6 hours via spbackup.</Description>
    <URI>\SPBackup-Library</URI>
  </RegistrationInfo>

  <Triggers>
    <CalendarTrigger>
      <Repetition>
        <Interval>PT6H</Interval>
        <StopAtDurationEnd>false</StopAtDurationEnd>
      </Repetition>
      <StartBoundary>2026-01-01T02:30:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>

  <Principals>
    <Principal id="Author">
      <UserId>S-1-5-18</UserId>          <!-- SYSTEM account; change if needed -->
      <LogonType>Password</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>

  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>false</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>false</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>true</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <!-- Library backups can be large — allow up to 6 hours -->
    <ExecutionTimeLimit>PT6H</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>

  <Actions Context="Author">
    <Exec>
      <Command>C:\Program Files\PowerShell\7\pwsh.exe</Command>
      <Arguments>-NoProfile -NonInteractive -File "C:\spbackup\spbackup.ps1" library backup --url "https://contoso.sharepoint.com/sites/team" --library "Documents" --out "C:\backups\docs" --verbose</Arguments>
      <WorkingDirectory>C:\spbackup</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
